# initial-setup.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/X86/initial-setup.sh

[Install]
WantedBy=default.target